# blockwatch.service
[Unit]
Description=Miningcore Blockwatch (writes blockfound.txt + blockfound.log)
After=network.target docker.service
Wants=docker.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=umbrel
Environment=MININGCORE_API=http://127.0.0.1:4000
Environment=POOL_ID=dgb
Environment=POLL_SEC=15
ExecStart=/usr/bin/python3 -u /home/umbrel/umbrel/app-data/willitmod-dev-dgb/data/pool/config/walletguard/scripts/blockwatch.py
Restart=on-failure
RestartSec=10
TimeoutStartSec=30

# Stop cleanly (docker logs -f can hang otherwise)
KillMode=control-group
TimeoutStopSec=10
KillSignal=SIGTERM
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
