# blockwatch.service
[Unit]
Description=Miningcore Blockwatch (writes blockfound.txt, bf_log.txt, height.txt, poolblock.txt)
After=network.target docker.service
Wants=docker.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=umbrel
WorkingDirectory=/home/umbrel/umbrel/app-data/NOC-Mining/scripts/blockwatch
ExecStart=/usr/bin/python3 -u /home/umbrel/umbrel/app-data/NOC-Mining/scripts/blockwatch/blockwatch.py
Restart=always
RestartSec=10
TimeoutStartSec=30

KillMode=control-group
TimeoutStopSec=10
KillSignal=SIGTERM
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
